# timetree-exporter/crontab
# crontab
# execute the exporter every hour
0 0 * * * * *    /usr/local/bin/execute-exporter.sh
